# .config/autostart/greenclip.desktop
[Desktop Entry]
Name=Greenclip
Exec=greenclip daemon
Type=Application
